# user/systemd/.config/systemd/user/clear-cache.timer
[Unit]
Description=clear various caches

[Timer] 
Persistent=true
OnCalendar=weekly
 
[Install]
WantedBy=timers.target 
