# user/systemd/.config/systemd/user/rbw-sync.timer
[Unit]
Description=rbw

[Timer] 
Persistent=true
OnBootSec=1min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target 
