# user/systemd/.config/systemd/user/vdirsyncer-meta.service
[Unit]
Description=vdirsyncer metasync

[Service]
ExecStart=/usr/bin/vdirsyncer metasync --max-workers=1 
RuntimeMaxSec=600
