# user/systemd/.config/systemd/user/vdirsyncer-sync.service
[Unit]
Description=vdirsyncer 

[Service]
ExecStart=/usr/bin/vdirsyncer sync --max-workers=1 
RuntimeMaxSec=600
