# systemd-user/.config/systemd/user/clear-cache.service
[Unit]
Description=clear various cache locations

[Service]
ExecStart=rm /home/ethan/.cache/screenshot/*
RuntimeMaxSec=60
