# systemd-user/.config/systemd/user/rbw-sync.service
[Unit]
Description=vdirsyncer metasync

[Service]
ExecStart=/usr/bin/rbw sync 
RuntimeMaxSec=600
